# periodic/daily/backup-database.sh
#!/bin/bash
/usr/bin/python3 /app/tasks_sql.py --interval=daily
