# periodic/weekly/backup-database.sh
python3 /app/tasks-sql.py --interval=weekly
